# saltstack/salt/dev/supervisor/mailcatcher.conf
[program:mailcatcher]
command = mailcatcher -f --http-ip 0.0.0.0
autostart = true
autorestart = true
stdout_logfile = /var/log/mailcatcher.stdout.log
stderr_logfile = /var/www/mailcatcher.stderr.log
